# build-arm64/_CPack_Packages/Linux/DEB/textsnippets_1.0.0_arm64/etc/xdg/autostart/textsnippets-autostart.desktop
[Desktop Entry]
Name=Textsnippets
Comment=Text snippet manager with system tray
Exec=textsnippets
Icon=accessories-text-editor
Terminal=false
Type=Application
Categories=Utility;TextTools;GTK;
X-GNOME-Autostart-enabled=true
StartupNotify=false
NoDisplay=true
